# chef/guckmon/files/default/chef.conf
description "Starts Chef Client on bootup with a 1800 second interval"
author "Ainsey11"
start on runlevel [2345]
exec chef-client -i 1800
